# packages/nostr-relay/nostr_relay-1.13.tar.gz/nostr_relay-1.13/config/nostr-relay.service
[Unit]
Description=uvicorn daemon for nostr relay
After=network.target

[Service]
Type=notify
# the specific user that our service will run as
User=www-data
Group=www-data
RuntimeDirectory=uvicorn
WorkingDirectory=/home/www/nostr
ExecStart=/home/www/envs/nostr/bin/nostr-relay -c /etc/nostr-relay.yaml serve
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
